# scd/systemd/cml-scd.socket
[Unit]
Description = GyroidOS CML-SCD Socket

[Socket]
FileDescriptorName=scd_control
Accept=no
ListenSequentialPacket=/run/socket/cml-scd_control
